# deploy/bos3000.service
[Unit]
Description=BOS3000 VoIP Callback Platform
After=network.target postgresql.service redis.service freeswitch.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=bos3000
Group=bos3000
WorkingDirectory=/opt/bos3000
EnvironmentFile=/opt/bos3000/.env
ExecStart=/opt/bos3000/bin/bos3000
Restart=on-failure
RestartSec=5
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bos3000

[Install]
WantedBy=multi-user.target
